# templates/etc/systemd/carbon-cache@.service.erb
[Unit]
Description=Graphite Carbon Cache instance %i
After=network.target

[Service]
Type=forking
PIDFile=<%= scope.lookupvar('graphite::gr_pid_dir') %>/carbon-cache-%i.pid
ExecStartPre=/bin/rm -f <%= scope.lookupvar('graphite::gr_pid_dir') %>/carbon-cache-%i.pid
ExecStart=/usr/bin/carbon-cache --config=<%= scope.lookupvar('graphite::config::carbon_conf_file') %> --instance=%i --pidfile=carbon-cache-%i.pid start
LimitNOFILE= <%= scope.lookupvar('graphite::gr_carbon_cache_ulimit') %>

[Install]
WantedBy=multi-user.target
